# data/pantheon-parental-controls.service.in
[Unit]
Description=Pantheon Digital Wellbeing Daemon

[Service]
Type=dbus
Restart=always
BusName=org.pantheon.ParentalControls
ExecStart=@DAEMON_PATH@

[Install]
WantedBy=multi-user.target
